# config/bot.conf
[program:<имя проекта>]
command=<путь до исполняемого python файла>
directory=<путь до директории проекта>
user=<имя пользователя хоста>
autorestart=true
redirect_stderr=true